# scripts/led_mmc/ledmmc.service
[Unit]
Description=MMC Status LED

[Service]
Type=simple
ExecStart=/usr/local/bin/ledmmc_enable.sh

[Install]
WantedBy=multi-user.target