# src/base_rootfs/etc/systemd/system/depmod_boot.service
[Unit]
Description=Depmod on boot

[Service]
Type=simple
ExecStart=/usr/local/bin/depmod_enable.sh

[Install]
WantedBy=multi-user.target