# src/base_rootfs/etc/systemd/system/display-sleep.service
[Unit]

Description=Display Sleep

After=default.target

[Service]

ExecStart=/usr/local/bin/display-sleep.sh

[Install]

WantedBy=default.target
