# scripts/usb_host/usbhost.service
[Unit]
Description=USB Host Service

[Service]
Type=simple
ExecStart=/usr/local/bin/usbhost_enable.sh

[Install]
WantedBy=multi-user.target